# cgv-alert/cgv-checker.service
[Unit]
Description=CGV IMAX schedule checker (continuous loop)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kelvinpark_dev
WorkingDirectory=/home/kelvinpark_dev/cgv-alert

# Chrome 경로 (기존 cron 환경변수와 동일)
Environment="PUPPETEER_EXECUTABLE_PATH=/usr/bin/google-chrome-stable"

# 디스코드 웹훅 — 배포 시 실제 URL로 교체하거나 EnvironmentFile 사용
Environment="DISCORD_WEBHOOK_URL=REPLACE_ME"

# ntfy 토픽 오버라이드가 필요하면:
# Environment="NTFY_TOPIC=kelvin-cgv-imax-v2"

ExecStart=/usr/bin/node checker.js

# 프로세스 죽으면 10초 후 자동 재시작
Restart=always
RestartSec=10

# 로그 (journalctl -u cgv-checker 로도 조회 가능)
StandardOutput=append:/home/kelvinpark_dev/cgv-alert/service.log
StandardError=append:/home/kelvinpark_dev/cgv-alert/service.log

[Install]
WantedBy=multi-user.target
